# ops/mcp.service
[Unit]
Description=XYTE MCP Service
After=network.target redis.service postgresql.service

[Service]
User=mcp
WorkingDirectory=/opt/mcp
Environment="REDIS_URL=redis://127.0.0.1:6379/0"
Environment="DATABASE_URL=postgresql+asyncpg://mcp:pass@127.0.0.1/mcp"
ExecStart=/opt/mcp-venv/bin/uvicorn xyte_mcp.server:app --host 0.0.0.0 --port 8000
Restart=always
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
